# alpine/openrc/mjpg-streamer
#!/sbin/openrc-run
# OpenRC service for mjpg-streamer (USB webcam)

name="mjpg-streamer"
description="mjpg-streamer HTTP service for USB webcam"

: ${USER:="video"}
: ${GROUP:="video"}
: ${LOG_DIR:="/var/log/mjpg-streamer"}
: ${PIDFILE:="/run/mjpg-streamer.pid"}
: ${PORT:=8090}
: ${DEVICE:="/dev/video0"}
: ${RESOLUTION:="1280x720"}
: ${FRAMERATE:=15}

# Plugin paths may vary by Alpine package version
: ${MJPG_BIN:="/usr/bin/mjpg_streamer"}
: ${INPUT_PLUGIN:="/usr/lib/mjpg-streamer/input_uvc.so"}
: ${OUTPUT_PLUGIN:="/usr/lib/mjpg-streamer/output_http.so"}
: ${WWW_DIR:="/usr/share/mjpg-streamer/www"}

command="${MJPG_BIN}"
command_args="-i \"${INPUT_PLUGIN} -d ${DEVICE} -r ${RESOLUTION} -f ${FRAMERATE}\" -o \"${OUTPUT_PLUGIN} -w ${WWW_DIR} -p ${PORT} -l 127.0.0.1\""
start_stop_daemon_args="--background --make-pidfile --pidfile ${PIDFILE} --user ${USER}:${GROUP}"

depend() {
    need localmount net
}

start_pre() {
    checkpath -d -o ${USER}:${GROUP} -m 0755 ${LOG_DIR}
}

start() {
    ebegin "Starting mjpg-streamer on 127.0.0.1:${PORT}"
    start-stop-daemon --start ${start_stop_daemon_args} \
        --exec /bin/sh -- -c "exec ${command} ${command_args} > ${LOG_DIR}/mjpg-streamer.log 2>&1"
    eend $?
}

stop() {
    ebegin "Stopping mjpg-streamer"
    start-stop-daemon --stop --pidfile ${PIDFILE} --retry 5
    eend $?
}
